# centos83-shadowsocks-privoxy/supervisor.d/default.ini
[program:shadowsocks]
command=sslocal -s %(ENV_S_SERVER_IP)s -p %(ENV_S_SERVER_PORT)s -l 1080 -k %(ENV_S_PASSWORD)s -t 600 -m %(ENV_S_METHOD)s
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true

[program:privoxy]
command=privoxy /etc/privoxy/config
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true